# deploy/crontab
# arxiv-context-feed cron jobs for Oracle Cloud
# Install with: crontab deploy/crontab

# Daily pipeline at 06:00 UTC (analyze + ingest + reddit + audit log — all on Oracle Cloud)
0 6 * * * cd /home/ubuntu/arxiv-context-feed && .venv/bin/contextual-arxiv-feed run-analyze -o /home/ubuntu/arxiv-feed-logs/decisions-$(date +\%Y\%m\%d).json >> /home/ubuntu/arxiv-feed-logs/analyze-$(date +\%Y\%m\%d).log 2>&1

# Weekly updates check (Sunday 08:00 UTC) - no LLM needed
0 8 * * 0 cd /home/ubuntu/arxiv-context-feed && .venv/bin/contextual-arxiv-feed run-updates --lookback-days 7 >> /home/ubuntu/arxiv-feed-logs/updates-$(date +\%Y\%m\%d).log 2>&1

# Weekly citations refresh (Saturday 10:00 UTC) - API calls only
0 10 * * 6 cd /home/ubuntu/arxiv-context-feed && .venv/bin/contextual-arxiv-feed refresh-citations >> /home/ubuntu/arxiv-feed-logs/citations-$(date +\%Y\%m\%d).log 2>&1

# Monthly ChromaDB prune (1st of month, 04:00 UTC) - delete papers older than 9 months
0 4 1 * * cd /home/ubuntu/arxiv-context-feed && .venv/bin/contextual-arxiv-feed prune-chromadb >> /home/ubuntu/arxiv-feed-logs/prune-$(date +\%Y\%m\%d).log 2>&1

# Log rotation: delete logs older than 30 days
0 0 * * * find /home/ubuntu/arxiv-feed-logs -name "*.log" -mtime +30 -delete 2>/dev/null
0 0 * * * find /home/ubuntu/arxiv-feed-logs -name "*.json" -mtime +30 -delete 2>/dev/null
