# deploy/llama-server.service
[Unit]
Description=llama.cpp Server for arXiv Feed Judge
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/arxiv-context-feed
ExecStart=/home/ubuntu/arxiv-context-feed/.venv/bin/python -m llama_cpp.server \
    --model /home/ubuntu/.cache/arxiv-feed/models/qwen2.5-14b-instruct-q4_k_m.gguf \
    --host 127.0.0.1 \
    --port 8080 \
    --n_ctx 4096 \
    --n_threads 4 \
    --n_gpu_layers 0
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
